# 99-asus-screenpad-permissions.rules
SUBSYSTEM=="module", KERNEL=="asus_nb_wmi", ACTION=="add", RUN+="/usr/bin/chmod ugo+w /sys/class/leds/asus::screenpad/brightness"
